# data/session-snapshot-20260327/jarvis-maintenance.service
[Unit]
Description=JARVIS Weekly Maintenance

[Service]
Type=oneshot
ExecStart=/usr/local/bin/jarvis-maintenance.sh
